# cookbooks/conky/templates/debian/ip.cron
0 */6 * * * <%= @user %> /usr/bin/dig +short myip.opendns.com @resolver1.opendns.com > /home/<%= @user =>/.config/conky/output/ip
