# cookbooks/conky/templates/debian/ip.cron
0 */6 * * * <%= @user %> curl ifconfig.me/ip > ~/.config/conky/output/ip