# cookbooks/conky/templates/debian/ip.cron
0 */6 * * * dan curl ifconfig.me/ip > /home/dan/.config/conky/output/ip